# files/etc/init.d/hostname
#!/sbin/openrc-run

description="Sets the hostname of the machine."

depend() {
        keyword -prefix -lxc -docker
}

start() {
        KOPT_hostname=$( sed "s#.*:\(p.*\):usb0.10.*#\1#" /proc/cmdline )
        if [ -n $KOPT_hostname ] ; then
                opts="$KOPT_hostname"
        elif [ -s /etc/hostname ] ; then
                opts="-F /etc/hostname"
        else
                opts="${hostname:-localhost}"
        fi
        ebegin "Setting hostname"
        hostname $opts
        eend $?
}

